# deploy/systemd/pmtool.service
[Unit]
Description=Projektmanager Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/florian/Projektmanager/Project
Environment=PMTOOL_DATA_DIR=/home/florian/Projektmanager/Project
Environment=PMTOOL_PUBLIC_BASE_URL=http://100.80.250.84:8765
Environment=PMTOOL_TRUST_PROXY_HEADERS=0
ExecStart=/home/florian/Projektmanager/Project/.venv/bin/python -m pmtool.collab_server
Restart=always
RestartSec=5
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
